# services/inotify/run
#!/bin/sh

while inotifywait -r -e modify /config/apache/site-confs; do
    /usr/sbin/apache2 -f /etc/apache/apache2.conf
done
